# polymarket-trading/pumaclaw-report.timer
[Unit]
Description=Run PumaClaw Daily Report at 08:00 UTC

[Timer]
OnCalendar=*-*-* 08:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
