# polymarket-trading/pumaclaw-straddle.service
# Instalación recomendada (EC2, usuario ubuntu):
#   cp polymarket-trading/pumaclaw-straddle.service ~/.config/systemd/user/
#   systemctl --user daemon-reload && systemctl --user enable --now pumaclaw-straddle.service
# Para wallet dedicada straddle, cambia POLY_ENV_FILE a ~/.openclaw/.env_straddle
#
# Si usas systemd a nivel sistema (sudo), añade User=ubuntu y Group=ubuntu bajo [Service].

[Unit]
Description=PumaClaw Double-Cheap Straddle Bot (Polymarket)
After=network-online.target

[Service]
Type=simple
Environment=HOME=/home/ubuntu
Environment=PYTHONUNBUFFERED=1
Environment=POLY_ENV_FILE=/home/ubuntu/.openclaw/.env
Environment=STRADDLE_USD_PER_LEG=2
Environment=STRADDLE_MAX_SHARES_PER_LEG=6
Environment=STRADDLE_LIMIT=0.30
Environment=STRADDLE_OTHER_WITHIN=0.01
Environment=STRADDLE_CONFIRM_SEC=45
Environment=STRADDLE_TIMEOUT_SEC=45
Environment=STRADDLE_STRATEGY_NAME=STRADDLE_O2_OPT_v1
Environment=STRADDLE_TRADES_LOG_FILE=/home/ubuntu/trades_history_straddle.json
Environment=STRADDLE_TRADED_MARKETS_FILE=/home/ubuntu/traded_markets_straddle.json
Environment=STRADDLE_LOCK_FILE=/home/ubuntu/straddle.lock
WorkingDirectory=/home/ubuntu/.openclaw/workspace/skills/polymarket
ExecStart=/home/ubuntu/.openclaw/workspace/skills/polymarket/.venv/bin/python3 /home/ubuntu/.openclaw/workspace/skills/polymarket/scripts/doublecheap_straddle_bot.py
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/straddle_bot.log
StandardError=append:/home/ubuntu/straddle_bot.log

[Install]
WantedBy=default.target
